# etc/supervisord.d/kopf_update.ini
[program:kopf-update]
command=/opt/qnib/kopf/bin/update.sh
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
