# etc/supervisord.d/kopf_update.ini
[program:kopf-update]
command=consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/nginx.conf.ctmpl:/etc/nginx/nginx.conf:nginx -s reload"
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
